# deploy/lhpc-web.service
# LoRaHAM Pi Control web console — CANONICAL managed unit (generated; do not hand-edit).
# Rendered by lhpc.core.updater_units; the one-click updater proves this file BYTE-EXACT.
# Restore with `lhpc self-update --repair-integration` after any change.
[Unit]
Description=LoRaHAM Pi Control web console (loopback-only)
Documentation=file://%h/loraham-pi-control/src/loraham-pi-control/docs/deployment.md
After=network-online.target
Wants=network-online.target lhpc-selfupdate.path lhpc-nginx-restart.path
# Refuse to (re)start mid-uninstall — the updater's OnFailure= must not resurrect the console
# while a teardown is removing this deployment.
ConditionPathExists=!%h/loraham-pi-control/.lhpc-uninstalling
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# INTENTIONAL: kill only the main console process, never its control group. LHPC identity-tracks and
# lifecycle-manages the LoRaHAM stacks and detached build/test jobs it starts; a web restart or a
# self-update must NOT tear those workloads down. Controller uninstall is the ONE path that stops and
# verifies them explicitly. The default KillMode=control-group would kill them on every web restart.
KillMode=process
Environment=LHPC_RUNTIME_ROOT=%h/loraham-pi-control
WorkingDirectory=%h/loraham-pi-control/src/loraham-pi-control
ExecStart=%h/loraham-pi-control/venv/lhpc/bin/lhpc web --socket
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/loraham-pi-control/logs/lhpc-web.log
StandardError=append:%h/loraham-pi-control/logs/lhpc-web.log
SyslogIdentifier=lhpc-web
# The console pushes its reachable URL as the unit's Status: line (sd_notify STATUS=).
# Type stays `simple` -> NO readiness handshake: a missed notification can never fail or
# delay startup. %t/systemd/notify is NOT in InaccessiblePaths below, so the sandbox allows it.
NotifyAccess=main
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# runtime root + /tmp for lhpc itself; %h/.meshcore_nm because a stack GUI started under this
# unit (meshcore-nodegui) keeps its sessions/settings/favourites there. %h expanded by systemd.
ReadWritePaths=%h/loraham-pi-control -%h/.meshcore_nm /tmp
Environment=PLATFORMIO_CORE_DIR=%h/loraham-pi-control/build/tool-cache/platformio
Environment=IDF_TOOLS_PATH=%h/loraham-pi-control/build/tool-cache/espressif
Environment=XDG_CACHE_HOME=%h/loraham-pi-control/build/tool-cache/cache
Environment=PIP_CACHE_DIR=%h/loraham-pi-control/build/tool-cache/pip
# Close the user-systemd escape: checkout code runs in this process, so deny it the user bus.
# (systemctl --user / systemd-run --user / --machine=@.host all fail; AF_UNIX + journald stay.)
InaccessiblePaths=%t/bus %t/systemd/private
PrivateTmp=false
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK AF_BLUETOOTH
RestrictNamespaces=true
LockPersonality=true
# MemoryDenyWriteExecute is OMITTED here — QEMU's TCG JIT (meshcom) needs W+X. The updater
# helper (no QEMU) DOES set it.
SystemCallArchitectures=native

[Install]
WantedBy=default.target
