# deploy/lhpc-selfupdate.service
# LoRaHAM Pi Control self-update helper — CANONICAL managed unit (generated; do not hand-edit).
# Started ONLY by lhpc-selfupdate.path when the web writes a request marker. Parameter-free:
# the normal/overwrite mode is read from the claimed request. No systemctl calls anywhere.
[Unit]
Description=LoRaHAM Pi Control self-update (stop console, update, restart)
Documentation=file://%h/loraham-pi-control/src/loraham-pi-control/docs/deployment.md
# Only the .path may start this; a manual `systemctl start` is refused (protects the console
# from an errant Conflicts= stop).
RefuseManualStart=yes
# Declarative console stop/restart — the helper never calls systemctl. Conflicts+After stops
# the console (freeing its shared controller-runtime lock) BEFORE the update; OnSuccess/OnFailure
# brings it back on EVERY terminal state (success, failure, timeout, kill).
Conflicts=lhpc-web.service
After=lhpc-web.service
OnSuccess=lhpc-web.service
OnFailure=lhpc-web.service
# Belt-and-braces: never run without a request present.
ConditionPathExists=%h/loraham-pi-control/state/selfupdate.request

[Service]
Type=oneshot
Environment=LHPC_RUNTIME_ROOT=%h/loraham-pi-control
WorkingDirectory=%h/loraham-pi-control/src/loraham-pi-control
ExecStart=%h/loraham-pi-control/venv/lhpc/bin/lhpc self-update --run-service
TimeoutStartSec=900
StandardOutput=append:%h/loraham-pi-control/logs/lhpc-selfupdate.log
StandardError=append:%h/loraham-pi-control/logs/lhpc-selfupdate.log
SyslogIdentifier=lhpc-selfupdate
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/loraham-pi-control /tmp
Environment=PIP_CACHE_DIR=%h/loraham-pi-control/build/tool-cache/pip
InaccessiblePaths=%t/bus %t/systemd/private
PrivateTmp=false
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
# The helper runs only git/pip/CPython (no QEMU) — W^X memory is enforced.
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
